# supervisor_nginx_php-fpm.ini
[program:php-fpm]
command=php-fpm7 -F
autostart               = true
autorestart             = true
stdout_logfile          = /var/log/supervisor/programs_out.log
stderr_logfile          = /var/log/supervisor/programs_err.log

[program:nginx]
command=nginx -g "daemon off;"
autostart               = true
autorestart             = true
stdout_logfile          = /var/log/supervisor/programs_out.log
stderr_logfile          = /var/log/supervisor/programs_err.log  